# deployment/website.service
# symlinked from /etc/systemd/system
[Unit]
Description=Gunicorn instance to serve Flask API
After=network.target

[Service]
User=tcarbone
WorkingDirectory=/home/tcarbone/website-api
ExecStart=/home/tcarbone/website-api/.venv/bin/gunicorn -w 2 -b 127.0.0.1:5000 wsgi:app
Restart=always

[Install]
WantedBy=multi-user.target
